# vhdbuilder/parts/linux/cloud-init/artifacts/aks-logrotate.service
[Unit]
Description=runs the logrotate utility for log rotation with a custom configuration
[Service]
ExecStart=/usr/local/bin/logrotate.sh