# vhdbuilder/parts/linux/cloud-init/artifacts/cgroup-memory-telemetry.service
[Unit]
Description=Emit system cgroup memory telemetry

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/scripts/cgroup-memory-telemetry.sh