# vhdbuilder/parts/linux/cloud-init/artifacts/update_certs.path
[Unit]
Description=Monitor the cert directory for changes

[Path]
PathModified=/opt/certs
Unit=update_certs.service

[Install]
WantedBy=multi-user.target