# vhdbuilder/parts/k8s/cloud-init/artifacts/kubesystem.slice
[Unit]
Description=Kubernetes system service slice
Documentation=man:systemd.special(7)
DefaultDependencies=no
Before=slices.target
Requires=-.slice
After=-.slice
