# parts/linux/cloud-init/artifacts/ubuntu/containerd-nvidia-config.service
[Unit]
Description=A oneshot service that overwrites the containerd configuration on GB200 SKUs
After=multi-user.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/opt/azure/gb200-containerd-config.sh
RemainAfterExit=yes

[Install]
WantedBy=graphical.target