# parts/linux/cloud-init/artifacts/99-azure-network.rules
SUBSYSTEM=="net", ACTION=="add|change", ENV{ID_NET_NAME_SLOT}=="enP*", RUN+="/bin/bash -c \"echo %k > /etc/azure-network/nics-to-configure && systemctl start azure-network-config.service\""
