# parts/linux/cloud-init/artifacts/99-azure-network.rules
# AKS Network Tuning udev rule
# This rule configures network interface RX ring buffer settings
# Applies to network interfaces with slot names (ID_NET_NAME_SLOT) matching enP* pattern
SUBSYSTEM=="net", ACTION=="add|change", ENV{ID_NET_NAME_SLOT}=="enP*", RUN+="/opt/azure-network/configure-azure-network.sh %k"
