# parts/linux/cloud-init/artifacts/99-azure-network.rules
# AKS Network Tuning udev rule
# This rule sets the network interface RX ring buffer to 2048 for multi-core systems
# Applies to network interfaces with names matching enP* pattern
SUBSYSTEM=="net", ACTION=="add|change", ENV{ID_NET_NAME_SLOT}=="enP*", RUN+="/bin/bash -c 'if [ \"$(/usr/bin/nproc)\" -ge 4 ]; then /usr/sbin/ethtool -G %k rx 2048; fi'"
