# parts/linux/cloud-init/artifacts/99-azure-network.rules
# AKS Network Tuning udev rule
# This rule configures network interface RX ring buffer settings
# Applies to network interfaces with names matching enP* pattern
SUBSYSTEM=="net", ACTION=="add|change", ENV{ID_NET_NAME_SLOT}=="enP*", RUN+="/bin/bash -c 'echo %k > /etc/azure-network/nics-to-configure && /opt/azure-network/configure-azure-network.sh'"
