# parts/linux/cloud-init/artifacts/99-aks-nics-to-configure.rules
SUBSYSTEM=="net", ACTION=="add", ENV{ID_NET_NAME_SLOT}=="enP*", RUN+="/bin/bash -c \"echo %k >> /run/nics-to-configure\""
