# parts/linux/cloud-init/artifacts/aks-logrotate.service
[Unit]
Description=runs the logrotate utility for log rotation
[Service]
ExecStart=/usr/local/bin/aks-logrotate.sh