# parts/linux/cloud-init/artifacts/aks-log-collector.timer
[Unit]
Description=AKS Log Collector Timer

[Timer]
OnActiveSec=0m
OnBootSec=5min
OnUnitActiveSec=60m

[Install]
WantedBy=timers.target
