# parts/linux/cloud-init/artifacts/aks-node-controller.service
[Unit]
Description=Parse contract and run csecmd
ConditionPathExists=/opt/azure/containers/aks-node-controller-config.json
After=cloud-init.target
After=oem-cloudinit.service enable-oem-cloudinit.service
Wants=cloud-init.target

[Service]
Type=oneshot
ExecStart=/opt/azure/containers/aks-node-controller provision --provision-config=/opt/azure/containers/aks-node-controller-config.json
RemainAfterExit=No

[Install]
WantedBy=cloud-init.target
WantedBy=oem-cloudinit.service
