# parts/linux/cloud-init/artifacts/aks-node-controller.service
[Unit]
Description=AKS Node Controller - provision and run csecmd
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/azure/containers/aks-node-controller-wrapper.sh
RemainAfterExit=No

[Install]
WantedBy=multi-user.target
