# parts/linux/cloud-init/artifacts/aks-local-dns.service
[Unit]
Description=AKS Local DNS
Wants=network.target
After=network.target
Before=kubelet.service
Before=containerd.service

[Service]
Type=notify
NotifyAccess=all
WatchdogSec=10
Restart=always
KillMode=mixed
TimeoutStopSec=30
Slice=aks-local-dns.slice
EnvironmentFile=-/etc/default/aks-local-dns
# only start if the nodepool is tagged EnableAKSLocalDNS=true
ExecCondition=/bin/bash -c "jq -e '.compute.tagsList | map(select(.name | test(\"EnableAKSLocalDNS\"; \"i\")))[0].value // \"false\" | test(\"true\"; \"i\")' <(curl -fsSL -H 'Metadata: true' --noproxy '*' 'http://169.254.169.254/metadata/instance?api-version=2021-02-01')"
ExecStartPre=/opt/azure/aks-local-dns/prestart.sh
ExecStart=/opt/azure/aks-local-dns/aks-local-dns.sh

[Install]
WantedBy=multi-user.target
