# parts/linux/cloud-init/artifacts/aks-localdns-hosts-setup.timer
[Unit]
Description=Run AKS LocalDNS hosts setup periodically

[Timer]
# Run on boot, subject to the RandomizedDelaySec jitter below (0-5s after boot).
OnBootSec=0
# Refresh ~10-15s after each run completes (10s base plus 0-5s of RandomizedDelaySec
# jitter). AKS critical FQDN IPs can change due to load balancer rotation, Traffic
# Manager failover, or regional DNS updates. Frequent refreshes keep stale IPs from
# causing the hosts plugin to serve unreachable addresses for long.
OnUnitInactiveSec=10s
# De-synchronize nodes so refreshes do not all hit DNS at the same time. Adds 0-5s of
# jitter on top of the boot and post-run timings noted above.
RandomizedDelaySec=5s
# Timer accuracy (how much systemd can delay)
AccuracySec=1s

[Install]
WantedBy=timers.target
