# parts/linux/cloud-init/artifacts/aks-check-network.service
[Unit]
Description=AKS Check Network
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/azure/containers/aks-check-network.sh
