# parts/linux/cloud-init/artifacts/aks-logrotate.service
[Unit]
Description=runs the logrotate utility for log rotation with a custom configuration
[Service]
ExecStart=/opt/bin/logrotate.sh
