# parts/linux/cloud-init/artifacts/cgroup-pressure-telemetry.service
[Unit]
Description=Emit system cgroup pressure telemetry

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/scripts/cgroup-pressure-telemetry.sh