# parts/linux/cloud-init/artifacts/cgroup-pressure-telemetry.service
[Unit]
Description=Updates certificates copied from AKS DS

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/scripts/cgroup-pressure-telemetry.sh