# parts/linux/cloud-init/artifacts/ci-syslog.service
[Unit]
Description=CI syslog service

[Service]
Type=simple
ExecStart=/usr/local/bin/ci-syslog.sh

[Install]
WantedBy=multi-user.target