# parts/linux/cloud-init/artifacts/containerd.service
[Unit]
Description=containerd daemon
[Service]
ExecStart=/usr/bin/containerd
[Install]
WantedBy=multi-user.target
#EOF