# parts/linux/cloud-init/artifacts/iotrace.service
[Unit]
Description=Trace IO on node startup

[Service]
ExecStart=/opt/scripts/iotrace.sh

[Install]
WantedBy=multi-user.target

