# parts/linux/cloud-init/artifacts/localdns-exporter@.service
[Unit]
Description=Localdns Prometheus Exporter Per-Connection Worker
Documentation=https://github.com/Azure/AgentBaker

[Service]
Type=simple
ExecStart=/opt/azure/containers/localdns/localdns_exporter.sh
StandardInput=socket
StandardOutput=socket
StandardError=journal

# Security hardening - run as unprivileged dynamic user
# Exporter only reads systemd properties and world-readable .prom file
DynamicUser=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
PrivateMounts=yes

[Install]
WantedBy=multi-user.target
