# parts/linux/cloud-init/artifacts/localdns-exporter.socket
[Unit]
Description=Localdns Prometheus Exporter Socket
Documentation=https://github.com/Azure/AgentBaker

[Socket]
ListenStream=9353
# Accept=yes means "Spawn a new service instance for every connection"
# This enables on-demand execution: the exporter script only runs when VMAgent scrapes
Accept=yes

[Install]
WantedBy=sockets.target
