# parts/linux/cloud-init/artifacts/localdns-exporter@.service
[Unit]
Description=Localdns Prometheus Exporter Per-Connection Worker
Documentation=https://github.com/Azure/AgentBaker

[Service]
Type=simple
ExecStart=/opt/azure/containers/localdns/localdns_exporter.sh
StandardInput=socket
StandardOutput=socket
StandardError=journal

# Security hardening
User=root
Group=root
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/
ReadWritePaths=/run
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
