# parts/linux/cloud-init/artifacts/nbcparser.service
[Unit]
Description=Parse contract and run csecmd
After=cloud-init.target
Wants=nbcparser.target

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/azure/containers/nbcparser.sh
RemainAfterExit=no

[Install]
WantedBy=nbcparser.target