# parts/linux/cloud-init/artifacts/nbcparser.target
[Unit]
Description=Custom Target for Cloud Init Completion

[Install]
WantedBy=multi-user.target