# parts/linux/cloud-init/artifacts/secure-tls-bootstrap.service
[Unit]
Description=AKS Secure TLS Bootstrap Client
ConditionPathExists=/opt/azure/containers/secure-tls-bootstrap-retry-cap.sh
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Wrapped via secure-tls-bootstrap-retry-cap.sh to bound per-VM retries when
# the bootstrap client is stuck. See AB#38327355 for context: a single stuck
# VM previously generated thousands of failure events and dominated QoS.
ExecStart=/opt/azure/containers/secure-tls-bootstrap-retry-cap.sh \
    --verbose \
    --ensure-authorized \
    --next-proto=aks-tls-bootstrap \
    --kubeconfig=/var/lib/kubelet/kubeconfig \
    --cert-dir=/var/lib/kubelet/pki \
    --cluster-ca-file=/etc/kubernetes/certs/ca.crt \
    --log-file=/var/log/azure/aks/secure-tls-bootstrap.log \
    $BOOTSTRAP_FLAGS

[Install]
WantedBy=multi-user.target
