# parts/linux/cloud-init/artifacts/update_certs.path
[Unit]
Description=Monitor the cert directory for changes

[Path]
PathChanged=/opt/certs
Unit=update-certs.service

[Install]
WantedBy=multi-user.target