# parts/linux/cloud-init/artifacts/node-exporter/baseline/etc/systemd/system/node-exporter-restart.service
[Unit]
Description=Restart node-exporter when certificates change

[Service]
Type=OneShot
ExecStart=/bin/systemctl restart node-exporter.service
