# parts/linux/cloud-init/artifacts/node-exporter/baseline/etc/systemd/system/node-exporter-restart.service
[Service]
Type=OneShot
ExecStart=/bin/systemctl restart node-exporter.service
